# mystartupscript.conf
start on runlevel [2345]
stop on runlevel [!2345]

exec /path/to/bach.py